# rootfs/etc/services.d/mysql-slow-query-log/run
#!/usr/bin/execlineb -P
logutil-service /var/run/s6/mysql-slow-query-log-fifo /var/log/mysql-slow-query-logs
