# rootfs/etc/services.d/mysql-general-log/run
#!/usr/bin/execlineb -P
logutil-service /var/run/s6/mysql-general-log-fifo /var/log/mysql-general-logs
