# rootfs/etc/services.d/mysql-error-log/run
#!/usr/bin/execlineb -P
logutil-service /var/run/s6/mysql-error-log-fifo /var/log/mysql-error-logs
